# test/fixtures/strong-pm.service
# strong-pm
# Generated by strong-service-systemd on 2014-08-22T00:00:00.000Z

[Unit]
Description=Node application server

# To start/stop strong-pm service manually, run:
#
# systemctl start strong-pm
# systemctl stop strong-pm

[Service]
Type=simple
ExecStart=/custom/path/to/node /custom/path/to/my/script.js
WorkingDirectory=/

User=nobody
Group=nogroup

# Environment variables:
# no environment variables set

# Allow many incoming connections
LimitNOFILE=infinity
# Allow core dumps for debugging
LimitCORE=infinity

StandardInput=null
StandardOutput=syslog
StandardError=syslog
Restart=always

[Install]
WantedBy=multi-user.target
